# doppler-image-factory/ansible/roles/common/files/mysqldump.cron
## Crontab for take MySQL dumps

0 * * * * mysqldump doppler > /tmp/doppler.sql > /tmp/doppler.sql ; /usr/local/bin/aws s3 cp /tmp/doppler.sql s3://doppler-shopify/backup/$ENVIRONMENT/$(date "+\%Y")/$(date "+\%m")/$(date "+\%d")/

## End cron